# devstack/files/monasca-transform/monasca_transform_init.conf
description "Monasca Transform"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

limit nofile 32768 32768

expect daemon

exec /etc/monasca/transform/init/start-monasca-transform.sh