# devstack/files/spark/spark-worker.conf
description "Spark Worker"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

limit nofile 32768 32768

expect daemon

exec /etc/spark/init/start-spark-worker.sh